# profile/airootfs/etc/skel/.config/autostart/mullvad-vpn.desktop
[Desktop Entry]
Type=Application
Name=Mullvad VPN
Exec=mullvad-vpn
Icon=mullvad-vpn
X-KDE-autostart-phase=2
NoDisplay=true
